# crontab
# Reset the demo environment every hour
42 * * * * /app/bin/rake -C /app demo:reset
